# anyremoteWebinterface/upstart/anyremoteAnd2html.conf
# HCAN

description	"anyremoteAnd2html daemon"

start on started gdm
stop on stopped gdm

expect fork
respawn

script
  sh /usr/bin/anyremoteAnd2html.sh
end script
