# anyremoteWebinterface/upstart/hcanhidController.conf
# HCAN

description	"hcanhid-Controller daemon"

start on started gdm
stop on stopped gdm

expect fork
respawn

script
  sh /usr/bin/hcanhidController.sh
end script
