# debian/hcan/etc/init/hcand.conf
# HCAN

description	"hcan daemon"

start on started gdm
stop on stopped gdm

expect fork
respawn

exec /usr/sbin/hcand
